# oasis_visualization/config/systemd/dismiss-osd.timer
[Unit]
Description=Schedule daily OSD dismissal at noon

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true
AccuracySec=1min
Unit=dismiss-osd.service

[Install]
WantedBy=timers.target

# To install:
# sudo cp dismiss-osd.timer /etc/systemd/system/

# To enable the timer:
# sudo systemctl enable --now dismiss-osd.timer

# To check status:
# systemctl status dismiss-osd.timer
