# oasis_visualization/config/systemd/dismiss_osd.timer
[Unit]
Description=Dismiss power OSD daily at 3 AM

[Timer]
OnCalender=hourly
Persistent=False

[Install]
WantedBy=timers.target
