# oasis_drivers_py/config/systemd/dismiss_osd.timer
[Unit]
Description=Dismiss power OSD daily at 3 AM

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=False

[Install]
WantedBy=dismiss_osd.target
