# oasis_drivers_py/config/systemd/dismiss_osd.service
[Unit]
Description=Dismiss power OSD
After=basic.target
Wants=basic.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -x gsd-power
