# oasis_drivers_py/config/systemd/oasis_drivers.service
[Unit]

Description=OASIS drivers
After=network.target

[Service]

#
# User configuration
#
# TODO: Better way to define this
#

User=ubuntu

#
# ROS 2 configuration
#

# ROS 2 distro and middleware
Environment=ROS2_DISTRO=@ROS2_DISTRO@
Environment=RMW_IMPLEMENTATION=rmw_fastrtps_dynamic_cpp

# Remove ROS 2 log timestamps and logger names to keep lines short and readable
Environment=RCUTILS_CONSOLE_OUTPUT_FORMAT="[{severity}]: {message}"

#
# Path configuration
#

Environment=OASIS_ENV=/home/${USER}/Documents/ros-ws/src/oasis/ros-ws/oasis-${ROS2_DISTRO}/install/setup.bash
Environment=OPENCV_LIB_DIR=/home/${USER}/Documents/ros-ws/src/oasis/ros-ws/cv/install/lib

#
# Scheduling configuration
#

# Favor hardware drivers over CPU-heavy perception services under saturation.
Nice=-5
CPUWeight=200
IOSchedulingClass=best-effort
IOSchedulingPriority=2

# Allow optional low-priority real-time scheduling for the standalone BNO086
# process. This grants only CAP_SYS_NICE to the service and caps FIFO/RR
# priorities at 20 so host-specific drop-ins can test prefixes such as:
#   Environment="OASIS_BNO086_LAUNCH_PREFIX=chrt -f 10"
# Keep FIFO priorities low, around 10. Do not use high FIFO priorities here:
# if the driver ever spins or blocks badly, it could starve other critical work.
LimitRTPRIO=20
LimitMEMLOCK=infinity
AmbientCapabilities=CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_NICE

# Optional per-process prefix for the standalone BNO086 driver only. This does
# not apply to camera, Telemetrix, magnetometer, or other driver nodes.
# Default keeps the BNO process ahead of other normal-priority driver work.
# For real-time testing, override this in a host-specific drop-in with:
#   Environment="OASIS_BNO086_LAUNCH_PREFIX=chrt -f 10"
# Verify scheduling with:
#   ps -eo pid,ni,rtprio,cls,comm,args | grep bno086
#   chrt -p <bno_pid>
Environment="OASIS_BNO086_LAUNCH_PREFIX=chrt -f 10"

#
# Service definition
#

# Delay launch until the system has been up long enough for ROS nodes to connect
ExecStart=/bin/bash -c 'set +o nounset && export LD_LIBRARY_PATH="`eval echo "${OPENCV_LIB_DIR}"`" && source `eval echo "${OASIS_ENV}"` && "$(ros2 pkg prefix oasis_drivers_py)/lib/oasis_drivers_py/delay_launch.sh" && ros2 launch oasis_drivers_py drivers_launch.py'

Restart=on-failure

[Install]

WantedBy=multi-user.target

# To install:
# sudo cp oasis_drivers.service /etc/systemd/system/

# To enable at startup:
# sudo systemctl enable oasis_drivers.service
